# etc/supervisor.d/cowrielog.ini
[program:cowrielog]
command=python3.7 /usr/lib/python3.7/site-packages/honeypot/exportlog.py
autostart=false
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/cowrielog.log
loglevel=warn
