# etc/supervisor.d/tornado.ini
[program:tornado]
command=python3.7 /home/honeypot/webserver/server.py
directory=/home/honeypot/webserver
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/tornado.log
loglevel=debug
